# infra/vm/systemd/aicore-searxng.service.example
[Unit]
Description=AI_ASSISTANT_CORE SearXNG runtime
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/home/bnhr/aicore/projects/core
ExecStart=/usr/bin/docker compose up -d searxng
ExecStop=/usr/bin/docker compose stop searxng
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target